# src/lib/systemd/system/olimex-ethernet-link-patch.service
[Unit]
Description=Olimex Ethernet Link Patch

[Service]
ExecStart=/usr/bin/olimex-el-patch
Type=oneshot

[Install]
WantedBy=multi-user.target